# ubuntu-18.04/network-lead/supervisor.dnsmasq.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq --conf-file=/config/dnsmasq.conf -k -d
redirect_stderr=true
autorestart=true
