# ubuntu-16.04/ubuntu-nodejs-dev-8/supervisor.nodejs.conf
[program:nodejs]
command=/usr/bin/npm start
directory=%(ENV_NODE_ROOT_DIRECTORY)s
autorestart=true